# _super_develop/go-push.ini
[program:go-push]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/go run /home/htdocs/gopush/main.go
autostart=true
autorestart=true
user=root
numprocs=1
redirect_stderr=true
stdout_logfile=/home/htdocs/gopush/logs/gopush.log
